# systemd/dabi-ingest-openintel-toplist.timer
[Unit]
Description=DABI OpenINTEL toplist forward-DNS ingest

[Timer]
OnCalendar=*-*-* 04:30:00 UTC
Persistent=true
RandomizedDelaySec=300
Unit=dabi-ingest-openintel-toplist.service

[Install]
WantedBy=timers.target
