# systemd/dabi-cert-renew.timer
[Unit]
Description=DABI weekly cert-renew timer

[Timer]
OnCalendar=Sun *-*-* 04:00:00 UTC
Persistent=true
RandomizedDelaySec=1800
Unit=dabi-cert-renew.service

[Install]
WantedBy=timers.target
