# systemd/dabi-ingest-openintel-ctlog.timer
[Unit]
Description=DABI OpenINTEL CT-log forward-DNS ingest

[Timer]
OnCalendar=*-*-* 21:30:00 UTC
Persistent=true
RandomizedDelaySec=300
Unit=dabi-ingest@openintel-ctlog.service

[Install]
WantedBy=timers.target
