# systemd/dabi-ingest-rir.timer
[Unit]
Description=DABI RIR-level WHOIS + rDNS daily

[Timer]
OnCalendar=*-*-* 23:30:00 UTC
Persistent=true
RandomizedDelaySec=300
Unit=dabi-ingest@rir.service

[Install]
WantedBy=timers.target
