# systemd/dabi-ingest-consolidate.timer
[Unit]
Description=DABI cross-dataset OpenSearch consolidation pass

[Timer]
OnCalendar=*-*-* 05:30:00 UTC
Persistent=true
RandomizedDelaySec=300
Unit=dabi-ingest@consolidate.service

[Install]
WantedBy=timers.target
