# systemd/dabi-ingest-openintel-zonefile.timer
[Unit]
Description=DABI OpenINTEL closed zone-based daily Parquet

[Timer]
OnCalendar=*-*-* 03:00:00 UTC   ; 08:30 IST
Persistent=true
RandomizedDelaySec=300
Unit=dabi-ingest@openintel-zonefile.service

[Install]
WantedBy=timers.target
