# systemd/dabi-ingest-openintel-cctld.timer
[Unit]
Description=DABI OpenINTEL ccTLD apex weekly ingest

[Timer]
OnCalendar=Mon *-*-* 23:00:00 UTC   ; 04:30 IST (Mondays)
Persistent=true
RandomizedDelaySec=300
Unit=dabi-ingest@openintel-cctld.service

[Install]
WantedBy=timers.target
