# systemd/dabi-ingest-czds.service
[Unit]
Description=DABI CZDS daily zone-file pull + index
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/srv/dabi/deploy
ExecStart=/srv/dabi/deploy/scripts/run-czds-pipeline.sh
TimeoutStartSec=4h
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=5
