# systemd/dabi-ingest-rir.timer
[Unit]
Description=DABI RIR ingest daily timer

[Timer]
OnCalendar=*-*-* 03:30:00 UTC
Persistent=true
RandomizedDelaySec=300
Unit=dabi-ingest-rir.service

[Install]
WantedBy=timers.target
