# systemd/dabi-ingest-tranco.service
[Unit]
Description=DABI Tranco daily top-1M pull + load
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/srv/dabi/deploy
ExecStart=/srv/dabi/deploy/scripts/run-tranco-pipeline.sh
TimeoutStartSec=30min
Nice=10
